# haproxy-config.upstart
# Upstart Configuration

description     "Golang Web App"
author          "Jeff Buchbinder"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
          
stop on runlevel [016]

respawn

pre-start script
        [ -r /etc/default/haproxy-config ] && . /etc/default/haproxy-config
	exec /usr/sbin/haproxy-config $HAPROXY_CONFIG_ARGS
end script

